# deploy/gunicorn.service
[Unit]
Description=FonoApp - Gunicorn Django
After=network.target postgresql.service

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/fonoapp
EnvironmentFile=/home/ubuntu/fonoapp/.env
ExecStart=/home/ubuntu/fonoapp/venv/bin/gunicorn \
    --workers 2 \
    --bind unix:/home/ubuntu/fonoapp/fonoapp.sock \
    fonoapp.wsgi:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
